# deploy/systemd/shadowdag-node.service
[Unit]
Description=ShadowDAG Node (testnet)
Documentation=https://github.com/ShadowDag/ShadowDag
After=network-online.target
Wants=network-online.target

# Start rate limiting: if the node crashes more than 5 times in 5 minutes
# systemd will give up. These directives must live in [Unit], not [Service]
# — older systemd versions warn (and newer ones reject) on the latter.
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=root
WorkingDirectory=/opt/ShadowDag
EnvironmentFile=/etc/shadowdag/node.env
ExecStart=/opt/ShadowDag/target/release/shadowdag-node --network=${NETWORK} --enable-stratum --enable-explorer --enable-ide

# Graceful shutdown — RocksDB must close cleanly
KillSignal=SIGTERM
TimeoutStopSec=90
SendSIGKILL=yes

# Auto-restart on crash (but not on clean exit)
Restart=on-failure
RestartSec=10

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Logs to journalctl
StandardOutput=journal
StandardError=journal
SyslogIdentifier=shadowdag-node

[Install]
WantedBy=multi-user.target
