# deploy/systemd/shadowdag-miner.service
[Unit]
Description=ShadowDAG Miner (testnet)
Documentation=https://github.com/ShadowDag/ShadowDag
After=shadowdag-node.service
Requires=shadowdag-node.service

# Start rate limiting belongs in [Unit], not [Service]. Older systemd
# versions warn (and newer ones reject) on the latter.
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=root
WorkingDirectory=/opt/ShadowDag
EnvironmentFile=/etc/shadowdag/node.env

# Wait up to 90s for the node's RPC port to accept connections.
ExecStartPre=/bin/bash /opt/ShadowDag/deploy/systemd/wait-for-rpc.sh ${RPC_PORT} 90

ExecStart=/opt/ShadowDag/target/release/shadowdag-miner --network=${NETWORK} --address=${MINER_ADDRESS} --threads=${MINER_THREADS} --rpc=127.0.0.1:${RPC_PORT}

KillSignal=SIGTERM
TimeoutStopSec=30

Restart=on-failure
RestartSec=15

LimitNOFILE=65536

StandardOutput=journal
StandardError=journal
SyslogIdentifier=shadowdag-miner

[Install]
WantedBy=multi-user.target
