# deploy/systemd/shadowdag-miner.service
[Unit]
Description=ShadowDAG Miner (testnet)
Documentation=https://github.com/ShadowDag/ShadowDag
After=shadowdag-node.service
Requires=shadowdag-node.service

# Start rate limiting belongs in [Unit], not [Service]. Older systemd
# versions warn (and newer ones reject) on the latter.
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=root
WorkingDirectory=/opt/ShadowDag
EnvironmentFile=/etc/shadowdag/miner.env

# Wait up to 90s for the node's RPC port to accept connections before
# starting the miner. See deploy/systemd/wait-for-rpc.sh for why this is
# a separate script instead of an inline curl (systemd quoting + curl -f
# don't play well together during node startup).
ExecStartPre=/bin/bash /opt/ShadowDag/deploy/systemd/wait-for-rpc.sh 19332 90

ExecStart=/opt/ShadowDag/target/release/shadowdag-miner --network=testnet --address=${MINER_ADDRESS} --threads=${MINER_THREADS} --rpc=127.0.0.1:19332

KillSignal=SIGTERM
TimeoutStopSec=30

Restart=on-failure
RestartSec=15

LimitNOFILE=65536

StandardOutput=journal
StandardError=journal
SyslogIdentifier=shadowdag-miner

[Install]
WantedBy=multi-user.target
